# ~/Library/LaunchAgents/com.ananda.log-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ananda.log-rotate</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/Michael/.pyenv/shims/python3</string>
        <string>/Users/Michael/Documents/development/mega-rag-chatbot/data_ingestion/crawler/log_rotate.py</string>
        <string>--log-dir</string>
        <string>/Users/Michael/Library/Logs/AnandaCrawler</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/Michael/Library/Logs/AnandaCrawler/log-rotate.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/Michael/Library/Logs/AnandaCrawler/log-rotate-error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>WorkingDirectory</key>
    <string>/Users/Michael/Documents/development/mega-rag-chatbot/data_ingestion/crawler</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/Michael/.pyenv/shims:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
