# data_ingestion/crawler/com.ananda.crawler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ananda.crawler.ananda-public</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Michael/.pyenv/versions/mega-rag-chatbot/bin/python</string>
        <string>/Users/Michael/Documents/development/mega-rag-chatbot/data_ingestion/crawler/crawler_supervisor.py</string>
        <string>--site</string>
        <string>ananda-public</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Michael/Documents/development/mega-rag-chatbot/data_ingestion/crawler</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/Michael/Documents/development/mega-rag-chatbot/data_ingestion</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/Michael/Library/Logs/AnandaCrawler/supervisor_ananda-public.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Michael/Library/Logs/AnandaCrawler/supervisor_ananda-public.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <!-- Resource limits -->
    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>10</integer>

    <!-- Restart on failure -->
    <key>OnDemand</key>
    <false/>

    <!-- Throttle restarts -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Don't start if system is sleeping -->
    <key>LaunchOnlyOnce</key>
    <false/>
</dict>
</plist>
